# deploy/systemd/ambient-snapshot.service
[Unit]
Description=Ambient Weather snapshot collector

[Service]
Type=oneshot
WorkingDirectory=/home/ben/py-projects/ambient-tool
ExecStart=/home/ben/py-projects/ambient-tool/.venv/bin/python -m ambient_tool.cli snapshot
